# crawler/crawler.conf
[supervisord]
nodaemon=true
environment=DISPLAY=":99"

[program:geckodriver-log-reader_1]
command=faust -A geckodriver_log_reader worker -p 6160
priority=100
stdout_logfile=logs/crawler_1/geckodriver_out.log
stderr_logfile=logs/crawler_1/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_1.json",
    GECKODRIVER_LOG_FILE="geckodriver_1.log"
[program:crawler_1]
command=faust -A crawler worker -p 6181
priority=200
environment=WS_PORT="6101"
stdout_logfile=logs/crawler_1/crawler1_out.log
stderr_logfile=logs/crawler_1/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_1.json"
[program:websocket_1]
command=faust -A websocket worker -p 6191
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_1.json",
    WS_PORT="6101"
stdout_logfile=logs/crawler_1/ws1_out.log
stderr_logfile=logs/crawler_1/ws1_err.log

[program:geckodriver-log-reader_2]
command=faust -A geckodriver_log_reader worker -p 6260
priority=100
stdout_logfile=logs/crawler_2/geckodriver_out.log
stderr_logfile=logs/crawler_2/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_2.json",
    GECKODRIVER_LOG_FILE="geckodriver_2.log"
[program:crawler_2]
command=faust -A crawler worker -p 6281
priority=200
environment=WS_PORT="6201"
stdout_logfile=logs/crawler_2/crawler1_out.log
stderr_logfile=logs/crawler_2/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_2.json"
[program:websocket_2]
command=faust -A websocket worker -p 6291
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_2.json",
    WS_PORT="6201"
stdout_logfile=logs/crawler_2/ws1_out.log
stderr_logfile=logs/crawler_2/ws1_err.log

[program:geckodriver-log-reader_3]
command=faust -A geckodriver_log_reader worker -p 6360
priority=100
stdout_logfile=logs/crawler_3/geckodriver_out.log
stderr_logfile=logs/crawler_3/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_3.json",
    GECKODRIVER_LOG_FILE="geckodriver_3.log"
[program:crawler_3]
command=faust -A crawler worker -p 6381
priority=200
environment=WS_PORT="6301"
stdout_logfile=logs/crawler_3/crawler1_out.log
stderr_logfile=logs/crawler_3/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_3.json"
[program:websocket_3]
command=faust -A websocket worker -p 6391
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_3.json",
    WS_PORT="6301"
stdout_logfile=logs/crawler_3/ws1_out.log
stderr_logfile=logs/crawler_3/ws1_err.log

[program:geckodriver-log-reader_4]
command=faust -A geckodriver_log_reader worker -p 6460
priority=100
stdout_logfile=logs/crawler_4/geckodriver_out.log
stderr_logfile=logs/crawler_4/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_4.json",
    GECKODRIVER_LOG_FILE="geckodriver_4.log"
[program:crawler_4]
command=faust -A crawler worker -p 6481
priority=200
environment=WS_PORT="6401"
stdout_logfile=logs/crawler_4/crawler1_out.log
stderr_logfile=logs/crawler_4/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_4.json"
[program:websocket_4]
command=faust -A websocket worker -p 6491
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_4.json",
    WS_PORT="6401"
stdout_logfile=logs/crawler_4/ws1_out.log
stderr_logfile=logs/crawler_4/ws1_err.log

[program:geckodriver-log-reader_5]
command=faust -A geckodriver_log_reader worker -p 6560
priority=100
stdout_logfile=logs/crawler_5/geckodriver_out.log
stderr_logfile=logs/crawler_5/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_5.json",
    GECKODRIVER_LOG_FILE="geckodriver_5.log"
[program:crawler_5]
command=faust -A crawler worker -p 6581
priority=200
environment=WS_PORT="6501"
stdout_logfile=logs/crawler_5/crawler1_out.log
stderr_logfile=logs/crawler_5/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_5.json"
[program:websocket_5]
command=faust -A websocket worker -p 6591
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_5.json",
    WS_PORT="6501"
stdout_logfile=logs/crawler_5/ws1_out.log
stderr_logfile=logs/crawler_5/ws1_err.log

[program:geckodriver-log-reader_6]
command=faust -A geckodriver_log_reader worker -p 6660
priority=100
stdout_logfile=logs/crawler_6/geckodriver_out.log
stderr_logfile=logs/crawler_6/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_6.json",
    GECKODRIVER_LOG_FILE="geckodriver_6.log"
[program:crawler_6]
command=faust -A crawler worker -p 6681
priority=200
environment=WS_PORT="6601"
stdout_logfile=logs/crawler_6/crawler1_out.log
stderr_logfile=logs/crawler_6/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_6.json"
[program:websocket_6]
command=faust -A websocket worker -p 6691
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_6.json",
    WS_PORT="6601"
stdout_logfile=logs/crawler_6/ws1_out.log
stderr_logfile=logs/crawler_6/ws1_err.log

[program:geckodriver-log-reader_7]
command=faust -A geckodriver_log_reader worker -p 6760
priority=100
stdout_logfile=logs/crawler_7/geckodriver_out.log
stderr_logfile=logs/crawler_7/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_7.json",
    GECKODRIVER_LOG_FILE="geckodriver_7.log"
[program:crawler_7]
command=faust -A crawler worker -p 6781
priority=200
environment=WS_PORT="6701"
stdout_logfile=logs/crawler_7/crawler1_out.log
stderr_logfile=logs/crawler_7/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_7.json"
[program:websocket_7]
command=faust -A websocket worker -p 6791
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_7.json",
    WS_PORT="6701"
stdout_logfile=logs/crawler_7/ws1_out.log
stderr_logfile=logs/crawler_7/ws1_err.log

[program:geckodriver-log-reader_8]
command=faust -A geckodriver_log_reader worker -p 6860
priority=100
stdout_logfile=logs/crawler_8/geckodriver_out.log
stderr_logfile=logs/crawler_8/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_8.json",
    GECKODRIVER_LOG_FILE="geckodriver_8.log"
[program:crawler_8]
command=faust -A crawler worker -p 6881
priority=200
environment=WS_PORT="6801"
stdout_logfile=logs/crawler_8/crawler1_out.log
stderr_logfile=logs/crawler_8/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_8.json"
[program:websocket_8]
command=faust -A websocket worker -p 6891
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_8.json",
    WS_PORT="6801"
stdout_logfile=logs/crawler_8/ws1_out.log
stderr_logfile=logs/crawler_8/ws1_err.log

[program:geckodriver-log-reader_9]
command=faust -A geckodriver_log_reader worker -p 6960
priority=100
stdout_logfile=logs/crawler_9/geckodriver_out.log
stderr_logfile=logs/crawler_9/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_9.json",
    GECKODRIVER_LOG_FILE="geckodriver_9.log"
[program:crawler_9]
command=faust -A crawler worker -p 6981
priority=200
environment=WS_PORT="6901"
stdout_logfile=logs/crawler_9/crawler1_out.log
stderr_logfile=logs/crawler_9/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_9.json"
[program:websocket_9]
command=faust -A websocket worker -p 6991
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_9.json",
    WS_PORT="6901"
stdout_logfile=logs/crawler_9/ws1_out.log
stderr_logfile=logs/crawler_9/ws1_err.log

[program:geckodriver-log-reader_10]
command=faust -A geckodriver_log_reader worker -p 61060
priority=100
stdout_logfile=logs/crawler_10/geckodriver_out.log
stderr_logfile=logs/crawler_10/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_10.json",
    GECKODRIVER_LOG_FILE="geckodriver_10.log"
[program:crawler_10]
command=faust -A crawler worker -p 61081
priority=200
environment=WS_PORT="61001"
stdout_logfile=logs/crawler_10/crawler1_out.log
stderr_logfile=logs/crawler_10/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_10.json"
[program:websocket_10]
command=faust -A websocket worker -p 61091
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_10.json",
    WS_PORT="61001"
stdout_logfile=logs/crawler_10/ws1_out.log
stderr_logfile=logs/crawler_10/ws1_err.log

[program:geckodriver-log-reader_11]
command=faust -A geckodriver_log_reader worker -p 61160
priority=100
stdout_logfile=logs/crawler_11/geckodriver_out.log
stderr_logfile=logs/crawler_11/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_11.json",
    GECKODRIVER_LOG_FILE="geckodriver_11.log"
[program:crawler_11]
command=faust -A crawler worker -p 61181
priority=200
environment=WS_PORT="61101"
stdout_logfile=logs/crawler_11/crawler1_out.log
stderr_logfile=logs/crawler_11/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_11.json"
[program:websocket_11]
command=faust -A websocket worker -p 61191
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_11.json",
    WS_PORT="61101"
stdout_logfile=logs/crawler_11/ws1_out.log
stderr_logfile=logs/crawler_11/ws1_err.log

[program:geckodriver-log-reader_12]
command=faust -A geckodriver_log_reader worker -p 61260
priority=100
stdout_logfile=logs/crawler_12/geckodriver_out.log
stderr_logfile=logs/crawler_12/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_12.json",
    GECKODRIVER_LOG_FILE="geckodriver_12.log"
[program:crawler_12]
command=faust -A crawler worker -p 61281
priority=200
environment=WS_PORT="61201"
stdout_logfile=logs/crawler_12/crawler1_out.log
stderr_logfile=logs/crawler_12/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_12.json"
[program:websocket_12]
command=faust -A websocket worker -p 61291
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_12.json",
    WS_PORT="61201"
stdout_logfile=logs/crawler_12/ws1_out.log
stderr_logfile=logs/crawler_12/ws1_err.log

[program:geckodriver-log-reader_13]
command=faust -A geckodriver_log_reader worker -p 61360
priority=100
stdout_logfile=logs/crawler_13/geckodriver_out.log
stderr_logfile=logs/crawler_13/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_13.json",
    GECKODRIVER_LOG_FILE="geckodriver_13.log"
[program:crawler_13]
command=faust -A crawler worker -p 61381
priority=200
environment=WS_PORT="61301"
stdout_logfile=logs/crawler_13/crawler1_out.log
stderr_logfile=logs/crawler_13/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_13.json"
[program:websocket_13]
command=faust -A websocket worker -p 61391
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_13.json",
    WS_PORT="61301"
stdout_logfile=logs/crawler_13/ws1_out.log
stderr_logfile=logs/crawler_13/ws1_err.log

[program:geckodriver-log-reader_14]
command=faust -A geckodriver_log_reader worker -p 61460
priority=100
stdout_logfile=logs/crawler_14/geckodriver_out.log
stderr_logfile=logs/crawler_14/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_14.json",
    GECKODRIVER_LOG_FILE="geckodriver_14.log"
[program:crawler_14]
command=faust -A crawler worker -p 61481
priority=200
environment=WS_PORT="61401"
stdout_logfile=logs/crawler_14/crawler1_out.log
stderr_logfile=logs/crawler_14/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_14.json"
[program:websocket_14]
command=faust -A websocket worker -p 61491
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_14.json",
    WS_PORT="61401"
stdout_logfile=logs/crawler_14/ws1_out.log
stderr_logfile=logs/crawler_14/ws1_err.log

[program:geckodriver-log-reader_15]
command=faust -A geckodriver_log_reader worker -p 61560
priority=100
stdout_logfile=logs/crawler_15/geckodriver_out.log
stderr_logfile=logs/crawler_15/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_15.json",
    GECKODRIVER_LOG_FILE="geckodriver_15.log"
[program:crawler_15]
command=faust -A crawler worker -p 61581
priority=200
environment=WS_PORT="61501"
stdout_logfile=logs/crawler_15/crawler1_out.log
stderr_logfile=logs/crawler_15/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_15.json"
[program:websocket_15]
command=faust -A websocket worker -p 61591
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_15.json",
    WS_PORT="61501"
stdout_logfile=logs/crawler_15/ws1_out.log
stderr_logfile=logs/crawler_15/ws1_err.log

[program:geckodriver-log-reader_16]
command=faust -A geckodriver_log_reader worker -p 61660
priority=100
stdout_logfile=logs/crawler_16/geckodriver_out.log
stderr_logfile=logs/crawler_16/geckodriver_err.log
environment=
    MANAGER_PARAMS_FILE="manager_params_16.json",
    GECKODRIVER_LOG_FILE="geckodriver_16.log"
[program:crawler_16]
command=faust -A crawler worker -p 61681
priority=200
environment=WS_PORT="61601"
stdout_logfile=logs/crawler_16/crawler1_out.log
stderr_logfile=logs/crawler_16/crawler1_err.log
environment=MANAGER_PARAMS_FILE="manager_params_16.json"
[program:websocket_16]
command=faust -A websocket worker -p 61691
priority=100
environment=
    MANAGER_PARAMS_FILE="manager_params_16.json",
    WS_PORT="61601"
stdout_logfile=logs/crawler_16/ws1_out.log
stderr_logfile=logs/crawler_16/ws1_err.log
