# config/datasaver.conf
[supervisord]
nodaemon=true
environment=
    PYTHONPATH="$PYTHONPATH:crawler",



## Parallel crawler 1

[program:datasaver-1]
command=faust -A crawler.datasaver.sqlite worker -p 6162
priority=100
stdout_logfile=logs/crawler_1/datasaver_out.log
stderr_logfile=logs/crawler_1/datasaver_err.log
environment=
    MANAGER_PARAMS_FILE="config/manager_params.json",
    CRAWL_NAME="gcp-parallel-whead-kafka-2-1",
    DATABASE_NAME="data/crawl-data-1.sqlite",

## Parallel crawler 2

[program:datasaver-2]
command=faust -A crawler.datasaver.sqlite worker -p 6262
priority=100
stdout_logfile=logs/crawler_2/datasaver_out.log
stderr_logfile=logs/crawler_2/datasaver_err.log
environment=
    MANAGER_PARAMS_FILE="config/manager_params.json",
    CRAWL_NAME="gcp-parallel-whead-kafka-2-2",
    DATABASE_NAME="data/crawl-data-2.sqlite",

## Parallel crawler 3

[program:datasaver-3]
command=faust -A crawler.datasaver.sqlite worker -p 6362
priority=100
stdout_logfile=logs/crawler_3/datasaver_out.log
stderr_logfile=logs/crawler_3/datasaver_err.log
environment=
    MANAGER_PARAMS_FILE="config/manager_params.json",
    CRAWL_NAME="gcp-parallel-whead-kafka-2-3",
    DATABASE_NAME="data/crawl-data-3.sqlite",

## Parallel crawler 4

[program:datasaver-4]
command=faust -A crawler.datasaver.sqlite worker -p 6462
priority=100
stdout_logfile=logs/crawler_4/datasaver_out.log
stderr_logfile=logs/crawler_4/datasaver_err.log
environment=
    MANAGER_PARAMS_FILE="config/manager_params.json",
    CRAWL_NAME="gcp-parallel-whead-kafka-2-4",
    DATABASE_NAME="data/crawl-data-4.sqlite",

## Parallel crawler 5

[program:datasaver-5]
command=faust -A crawler.datasaver.sqlite worker -p 6562
priority=100
stdout_logfile=logs/crawler_5/datasaver_out.log
stderr_logfile=logs/crawler_5/datasaver_err.log
environment=
    MANAGER_PARAMS_FILE="config/manager_params.json",
    CRAWL_NAME="gcp-parallel-whead-kafka-2-5",
    DATABASE_NAME="data/crawl-data-5.sqlite",

## Parallel crawler 6

[program:datasaver-6]
command=faust -A crawler.datasaver.sqlite worker -p 6662
priority=100
stdout_logfile=logs/crawler_6/datasaver_out.log
stderr_logfile=logs/crawler_6/datasaver_err.log
environment=
    MANAGER_PARAMS_FILE="config/manager_params.json",
    CRAWL_NAME="gcp-parallel-whead-kafka-2-6",
    DATABASE_NAME="data/crawl-data-6.sqlite",

## Parallel crawler 7

[program:datasaver-7]
command=faust -A crawler.datasaver.sqlite worker -p 6762
priority=100
stdout_logfile=logs/crawler_7/datasaver_out.log
stderr_logfile=logs/crawler_7/datasaver_err.log
environment=
    MANAGER_PARAMS_FILE="config/manager_params.json",
    CRAWL_NAME="gcp-parallel-whead-kafka-2-7",
    DATABASE_NAME="data/crawl-data-7.sqlite",

## Parallel crawler 8

[program:datasaver-8]
command=faust -A crawler.datasaver.sqlite worker -p 6862
priority=100
stdout_logfile=logs/crawler_8/datasaver_out.log
stderr_logfile=logs/crawler_8/datasaver_err.log
environment=
    MANAGER_PARAMS_FILE="config/manager_params.json",
    CRAWL_NAME="gcp-parallel-whead-kafka-2-8",
    DATABASE_NAME="data/crawl-data-8.sqlite",

## Parallel crawler 9

[program:datasaver-9]
command=faust -A crawler.datasaver.sqlite worker -p 6962
priority=100
stdout_logfile=logs/crawler_9/datasaver_out.log
stderr_logfile=logs/crawler_9/datasaver_err.log
environment=
    MANAGER_PARAMS_FILE="config/manager_params.json",
    CRAWL_NAME="gcp-parallel-whead-kafka-2-9",
    DATABASE_NAME="data/crawl-data-9.sqlite",

## Parallel crawler 10

[program:datasaver-10]
command=faust -A crawler.datasaver.sqlite worker -p 61062
priority=100
stdout_logfile=logs/crawler_10/datasaver_out.log
stderr_logfile=logs/crawler_10/datasaver_err.log
environment=
    MANAGER_PARAMS_FILE="config/manager_params.json",
    CRAWL_NAME="gcp-parallel-whead-kafka-2-10",
    DATABASE_NAME="data/crawl-data-10.sqlite",

## Parallel crawler 11

[program:datasaver-11]
command=faust -A crawler.datasaver.sqlite worker -p 61162
priority=100
stdout_logfile=logs/crawler_11/datasaver_out.log
stderr_logfile=logs/crawler_11/datasaver_err.log
environment=
    MANAGER_PARAMS_FILE="config/manager_params.json",
    CRAWL_NAME="gcp-parallel-whead-kafka-2-11",
    DATABASE_NAME="data/crawl-data-11.sqlite",

## Parallel crawler 12

[program:datasaver-12]
command=faust -A crawler.datasaver.sqlite worker -p 61262
priority=100
stdout_logfile=logs/crawler_12/datasaver_out.log
stderr_logfile=logs/crawler_12/datasaver_err.log
environment=
    MANAGER_PARAMS_FILE="config/manager_params.json",
    CRAWL_NAME="gcp-parallel-whead-kafka-2-12",
    DATABASE_NAME="data/crawl-data-12.sqlite",

