# kafka/kafka1.conf
[supervisord]
nodaemon=true

# Kafka

[program:kafka]
command=/home/sbird/kafka_2.12-2.3.0/bin/kafka-server-start.sh config/server1.properties
priority=20
