# kafka/kafka-manager.conf
[supervisord]
nodaemon=true

# Kafka

[program:kafka-manager]
command=/opt/kafka/kafka-manager-2.0.0.2/kafka-manager -Dkafka-manager.zkhosts="0.0.0.0:2181"
priority=10
