# kafka/kafka.conf
[supervisord]
nodaemon=true

# Kafka

[program:zookeeper]
command=/opt/kafka/kafka_2.12-2.3.0/bin/zookeeper-server-start.sh config/zookeeper.properties
priority=10
stdout_logfile=../logs/kafka/zookeeper_out.log
stderr_logfile=../logs/kafka/zookeeper_err.log

[program:kafka-1]
command=/opt/kafka/kafka_2.12-2.3.0/bin/kafka-server-start.sh config/server1.properties
priority=20
stdout_logfile=../logs/kafka/kafka1_out.log
stderr_logfile=../logs/kafka/kafka1_err.log

[program:kafka-2]
command=/opt/kafka/kafka_2.12-2.3.0/bin/kafka-server-start.sh config/server2.properties
priority=20
stdout_logfile=../logs/kafka/kafka2_out.log
stderr_logfile=../logs/kafka/kafka2_err.log

[program:kafka-3]
command=/opt/kafka/kafka_2.12-2.3.0/bin/kafka-server-start.sh config/server3.properties
priority=20
stdout_logfile=../logs/kafka/kafka3_out.log
stderr_logfile=../logs/kafka/kafka3_err.log
